# persistence/upstart/nimbus-miner.conf
description "Start Nimbus miner"

start on runlevel [2345]

script
	exec ${NIMBUS_MINER_PATH}/nimbuspool-client-linux-x64 --wallet-address=${WALLET_ADDRESS} --extra-data=${WORKER_ID}
end script
